# talkingrobot.service
[Unit]
Description=Push-to-talk voice assistant (Pi + ReSpeaker 2-Mics HAT)
After=network-online.target sound.target
Wants=network-online.target

[Service]
# Run as unprivileged user that belongs to 'audio' group
User=admin
WorkingDirectory=/home/admin/talkingrobot

# Use the project venv interpreter
ExecStart=/home/admin/talkingrobot/venv/bin/python /home/admin/talkingrobot/app.py

# App configuration (contains keys via env but NOT secrets in code)
EnvironmentFile=/home/admin/talkingrobot/.env

# Restart resilience
Restart=always
RestartSec=3

# A little hardening (optional; relax if needed)
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
